# installation/metrics-exporter.service
[Unit]
Description=RAG-Search-LAB Metrics Exporter
After=network.target

[Service]
Type=simple
User=${INSTALL_USER}
WorkingDirectory=${INSTALL_BASE}
EnvironmentFile=${INSTALL_BASE}/installation/config.env
ExecStart=${INSTALL_BASE}/installation/metrics-exporter.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
